# build_dir/target-mips_r2_uClibc-0.9.33.2/root-ar71xx/usr/lib/opkg/info/6relayd.list
/usr/sbin/6relayd
/etc/hotplug.d/iface/30-6relay
/etc/config/6relayd
/usr/sbin/6relayd-update
/etc/init.d/6relayd
